# assets/installation/ubuntu/files/etc/systemd/system/exasol_launcher_host_version_check.service
# Host-side launcher version check
#
# Purpose: Best-effort fallback scheduled check for launcher updates
# Runs as: ubuntu user
# Runs on: Access node only, and only when DB version checks are enabled
# Dependencies: Requires database readiness marker and network
#
# Failure semantics:
# - Must never impact deployment availability.
# - Script exits 0 even when request fails/timeouts.

[Unit]
Description=Exasol Launcher Host Version Check
Wants=network-online.target
After=network-online.target exasol_launcher_ready.service
ConditionPathExists=/var/lib/exasol_launcher/state/ready.complete
ConditionPathExists=/etc/exasol_launcher/installation.json

[Service]
Type=oneshot
ExecCondition=/opt/exasol_launcher/scripts/isAccessNode.sh
ExecCondition=/opt/exasol_launcher/scripts/isDbVersionCheckEnabled.sh
ExecStart=/opt/exasol_launcher/scripts/hostLauncherVersionCheck.sh
User=ubuntu
WorkingDirectory=~
