# assets/installation/ubuntu/files/etc/systemd/system/exasol_launcher_prepare.service
# System Preparation Phase
#
# Purpose: Prepares the system for Exasol installation
# Runs as: root (requires system configuration permissions)
# Dependencies: None (first phase in installation workflow)
#
# Tasks:
# - Downloads the c4 installer binary
# - Configures udev rules for EBS volume access
# - Runs c4 preplay to prepare system for rootless container execution
#
# Idempotency: Only runs if exasol_launcher_prepare.complete marker doesn't exist
# Progress logging: Logs phase start/completion to journal for launcher monitoring

[Unit]
Description=Exasol Launcher System Preparation
ConditionPathExists=!/var/lib/exasol_launcher/state/exasol_launcher_prepare.complete

[Service]
Type=oneshot
ExecStart=/opt/exasol_launcher/scripts/prepareExasol.sh
ExecStart=/opt/exasol_launcher/scripts/runInfraHookScripts.sh pre-install-root
ExecStartPost=/usr/bin/touch /var/lib/exasol_launcher/state/exasol_launcher_prepare.complete
WorkingDirectory=/root
