# tools/buildvm/cloud-init/scripts/import-shared-keys.initd
#!/sbin/openrc-run

name="import-shared-keys"
description="Import SSH keys from shared folder"

depend() {
  need localmount
  after networking
}

start() {
  ebegin "Importing SSH keys from shared folder"
  /usr/local/bin/import-shared-keys.sh
  eend $?
}
